# .config/autostart/lunarclient.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=lunarclient
Comment=lunarclientstartup script
Exec=/tmp/.mount_lunarcDlpOFg/lunarclient --hidden
StartupNotify=false
Terminal=false